# Payload.PlayStationClassic/project_eris/network/etc/udev/rules.d/40-bluetooth.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="systemctl restart bluetooth"
